# chat/docker-scrapalot/supervisord.conf
; ============================================================================
; Scrapalot Workers — supervisord-managed per-queue Celery processes.
;
; PRD: docs/prd-competitive/prd_onyx_worker_stability.md (Phase 3).
; Adapted from onyx-dot-app/onyx backend/supervisord.conf.
;
; Before Phase 3: one Celery worker process served BOTH `documents` and
; `fast` queues with shared concurrency=1 thread. A 36-min reprocess on
; documents queue starved Annas restores on fast queue. If the worker
; SIGSEGV/SIGKILL'd, beat scheduler died with it and nothing recovered.
;
; After Phase 3: three independent processes inside the same container.
; The documents-queue crash (e.g. unhandled exception in
; reprocess_document_task watchdog) is isolated — fast queue keeps
; serving, beat keeps ticking. supervisord auto-restarts the crashed
; program within 10s (`startsecs=10`).
;
; All three programs share the same image / src/ bind-mount / Python
; deps — they differ only in --queues and --concurrency.
;
; Notes:
;   * `--pool=threads` (Phase 0) — no prefork SIGSEGV.
;   * `--max-tasks-per-child` removed (Phase 0) — irrelevant for threads.
;   * `stopasgroup=true` — supervisord sends SIGTERM to the whole process
;     group so SubprocessJobClient spawn-children die with the parent
;     worker on `supervisorctl stop` / container SIGTERM. Without this
;     they'd survive as orphans (container `stop_grace_period: 120s`
;     wouldn't reap them).
;   * `autorestart=true startsecs=10` — restart on crash but only count
;     it as "successfully started" if it lives ≥10s (avoids fast restart
;     loops on hard config errors).
; ============================================================================

[supervisord]
nodaemon=true
; `user=` removed — entrypoint.sh already drops root → uid 1000 via setpriv
; before exec'ing supervisord, so the `user=` directive is a no-op (it only
; switches uid when supervisord runs as root). Logfile + pidfile paths
; under /tmp and /app/data/logs which are uid-1000-writable per entrypoint.sh
; bootstrap.
logfile=/app/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/app/data/logs

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ----------------------------------------------------------------------------
; Documents queue — heavy reprocess + initial process_document tasks.
; Concurrency=1 thread; Phase 1 SubprocessJobClient spawns its own
; subprocess per reprocess, so this slot is just the watchdog thread.
; ----------------------------------------------------------------------------
[program:celery_worker_documents]
command=celery -A src.main.workers.celery_app worker
    --hostname=documents@%%h
    --queues=documents
    --pool=threads
    --concurrency=1
    --loglevel=info
; Worker logs go to file only. supervisord (uid 1000) cannot open
; /dev/stdout or /dev/fd/1 for write in this container — both fail with
; EACCES at the `make dispatchers` step (the rotation-aware logfile
; opener inside supervisord rejects non-regular targets even with
; stdout_logfile_maxbytes=0). Live-verified 2026-05-23: switching to
; /dev/fd/1 put all 3 workers into FATAL. Reverted to file logging.
;
; Inspect logs via:
;   docker exec scrapalot-workers tail -f /app/data/logs/celery_worker_documents.log
;   docker exec scrapalot-workers supervisorctl -c /app/supervisord.conf tail -f celery_worker_documents
stdout_logfile=/app/data/logs/celery_worker_documents.log
stdout_logfile_maxbytes=16MB
stdout_logfile_backups=3
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
priority=20

; ----------------------------------------------------------------------------
; Fast queue — Annas restores, summary backfill, hierarchy rebuild.
; Concurrency=4 threads since these are network/DB-bound and lightweight.
; Threads barely cost memory (single Python process; +~50 MB per thread
; for stack + DB session).
; ----------------------------------------------------------------------------
[program:celery_worker_fast]
command=celery -A src.main.workers.celery_app worker
    --hostname=fast@%%h
    --queues=fast
    --pool=threads
    --concurrency=4
    --loglevel=info
stdout_logfile=/app/data/logs/celery_worker_fast.log
stdout_logfile_maxbytes=16MB
stdout_logfile_backups=3
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
priority=10

; ----------------------------------------------------------------------------
; Beat scheduler — single instance for the system. Survives worker
; crashes independently. Schedule state persisted to /tmp so a restart
; doesn't double-fire scheduled tasks.
; ----------------------------------------------------------------------------
[program:celery_beat]
command=celery -A src.main.workers.celery_app beat
    --loglevel=info
    --schedule=/tmp/celerybeat-schedule
stdout_logfile=/app/data/logs/celery_beat.log
stdout_logfile_maxbytes=16MB
stdout_logfile_backups=3
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
priority=30

; ----------------------------------------------------------------------------
; Beat watchdog — catches the failure mode supervisord's autorestart can't:
; beat process up but scheduler frozen (deadlock, asyncio wedge), or beat
; dispatching but fast-queue workers all hung. Polls a Redis key written
; every 60 s by `scrapalot.celery_beat_heartbeat` (dispatched by beat,
; executed by a worker on the fast queue). Two-gate restart: 5 consecutive
; missed polls AND key age ≥ 10 min. See script docstring for tunables.
; ----------------------------------------------------------------------------
[program:beat_watchdog]
command=python /app/scripts/supervisord_watchdog_beat.py
stdout_logfile=/app/data/logs/beat_watchdog.log
stdout_logfile_maxbytes=4MB
stdout_logfile_backups=3
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
priority=40

; [program:log-redirect] (tail -qF → /dev/stdout) was removed: supervisord
; (uid 1000) cannot open /dev/stdout for write in this container, both
; /dev/stdout and /dev/fd/1 fail with EACCES at supervisord's dispatcher
; setup. Worker logs persist to /app/data/logs/*.log (file-only). Docker
; logs scrapalot-workers only shows supervisord-level events; full worker
; output requires `docker exec scrapalot-workers tail -f /app/data/logs/*.log`
; or `supervisorctl tail -f <program>`.
